# step-ca-client/rootfs/etc/s6-overlay/s6-rc.d/renewal-daemon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: step-ca-client
# Runs renewal-daemon script
# ==============================================================================
set -e

sleep 2 # give time to logging to flush and not be mixed

source /usr/bin/helpers.sh
set_debug

bashio::log.info "Starting step-ca-client"
/usr/bin/root-ca.sh

CERTFILE="/ssl/$(bashio::config 'certfile')"
if ! $(
    if [[ ${STEPDEBUG} -eq 1 ]];then set -x; fi;
    step certificate verify \
        "${CERTFILE}" \
        -roots=/root/.step/certs/root_ca.crt
)
then
    /usr/bin/create-with-token.sh
fi

bashio::log.info "Starting renewal daemon..."
exec /usr/bin/renewal-daemon.sh
